# docs/systemd/agent-memory-semwatch.service
# systemd user unit — keep one project's semantic index hot with `semwatch`.
#
# Preferred install (per project, no root needed):
#   bash docs/systemd/install-semwatch.sh /path/to/project
#
# Manual install:
#   mkdir -p ~/.config/systemd/user
#   cp docs/systemd/agent-memory-semwatch.service ~/.config/systemd/user/agent-memory-semwatch@.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now agent-memory-semwatch@$(systemd-escape /path/to/project).service
#
# Logs:  journalctl --user -u agent-memory-semwatch@$(systemd-escape /path/to/project) -f
# Stop:  systemctl --user disable --now 'agent-memory-semwatch@*'

[Unit]
Description=agent-memory semwatch (%i)
# Ollama down is fine: the indexer degrades to BM25/keyword and fills vectors on the next pass.
After=default.target

[Service]
Type=simple
# User units do not inherit a login shell PATH — without ~/.local/bin the
# console script from `uv pip install --user` / hatch is invisible (exit 127).
Environment=PATH=%h/.local/bin:%h/.cargo/bin:/usr/local/bin:/usr/bin:/bin
# Python block-buffers stdout when it is not a TTY — without this, `semwatch:`
# build lines sit in the buffer and journalctl shows nothing but stderr.
Environment=PYTHONUNBUFFERED=1
# %I = the UNescaped instance name (the real project path). Use %I, not %i:
# %i leaves the path escaped (hyphens become \x2d), so --root would receive a
# mangled path. The instance NAME keeps the escaped form (valid filename); %I
# restores it for the argument.
# Prefer absolute binary when install-semwatch.sh writes a drop-in override.
ExecStart=/usr/bin/env agent-memory --root %I semwatch --interval 2 --debounce 1
Restart=on-failure
RestartSec=10
# Cheap by design (polling a few dozen markdown files), but keep it polite anyway.
Nice=10

[Install]
WantedBy=default.target
